# systemd/securiphant-check-nfc.service
[Unit]
Description=securiphant-alert-bot

[Service]
ExecStart=$HOME/.local/bin/securiphant-check-nfc

[Install]
WantedBy=default.target
