# systemd/securiphant-check-environment.service
[Unit]
Description=securiphant-check-environment

[Service]
ExecStart=%h/.local/bin/securiphant-check-environment

[Install]
WantedBy=default.target
