# systemd/securiphant-alert-bot.service
[Unit]
Description=securiphant-alert-bot

[Service]
ExecStart=$HOME/.local/bin/securiphant-alert-bot

[Install]
WantedBy=default.target
