# systemd/securiphant-nfc-sensor.service
[Unit]
Description=securiphant-check-nfc

[Service]
ExecStart=%h/.local/bin/securiphant nfc-sensor -q

[Install]
WantedBy=default.target
