# systemd/securiphant-door-sensor.service
[Unit]
Description=securiphant-check-door

[Service]
ExecStart=%h/.local/bin/securiphant door-sensor -q

[Install]
WantedBy=default.target
