# systemd/securiphant-check-nfc.service
[Unit]
Description=securiphant-check-nfc

[Service]
ExecStart=%h/.local/bin/securiphant-check-nfc

[Install]
WantedBy=default.target
