# systemd/securiphant-display.service
[Unit]
Description=securiphant-display

[Service]
ExecStart=%h/.local/bin/securiphant display

[Install]
WantedBy=default.target
