# systemd/securiphant-speaker.service
[Unit]
Description=securiphant-speaker

[Service]
ExecStart=%h/.local/bin/securiphant speaker -q

[Install]
WantedBy=default.target
