# systemd/securiphant-display.service
[Unit]
Description=securiphant-alert-bot

[Service]
ExecStart=%h/.local/bin/securiphant-display

[Install]
WantedBy=default.target
