# systemd/securiphant-check-door.service
[Unit]
Description=securiphant-check-door

[Service]
ExecStart=%h/.local/bin/securiphant-check-door

[Install]
WantedBy=default.target
