# systemd/securiphant-check-door.service
[Unit]
Description=securiphant-alert-bot

[Service]
ExecStart=$HOME/.local/bin/securiphant-check-door

[Install]
WantedBy=default.target
